# dot_config/autostart/private_executable_xbindkeys.desktop
[Desktop Entry]
Exec=xbindkeys
Icon=dialog-scripts
Name=xbindkeys.sh
Type=Application
X-KDE-AutostartScript=true
